# home/user/.config/autostart/GoogleDrive.desktop
[Desktop Entry]
Exec=google-drive-ocamlfuse ~/migoogledrive
Icon=google
Name=google-drive-ocamlfuse
Path=
Terminal=False
Type=Application
