# home/user/.config/autostart/JuK.desktop
[Desktop Entry]
Exec=juk
Icon=
Name=JuK
Path=
Terminal=False
Type=Application
