# home/user/.config/autostart/Muezzin.desktop
[Desktop Entry]
Exec=muezzin
Icon=muezzin
Name=Muezzin
Path=
Terminal=False
Type=Application
