# home/user/.config/autostart/Obsidian.desktop
[Desktop Entry]
Exec=obsidian
Icon=obsidian
Name=Obsidian
Path=
Terminal=False
Type=Application
