# home/user/.config/autostart/PSensor.desktop
[Desktop Entry]
Exec=psensor
Icon=psensor
Name=PSensor
Path=
Terminal=False
Type=Application
