# home/user/.config/autostart/StartNODE.desktop
[Desktop Entry]
Exec=sh /home/mustafa/.0/startkit/start.sh
Icon=
Name=StartNODE
Path=
Terminal=False
Type=Application
