# home/mustafa/.config/autostart/guake.desktop
[Desktop Entry]
Exec=guake
Icon=guake
Name=Guake
Path=
Terminal=False
Type=Application
