# home/mustafa/.config/autostart/StartMYBASH.desktop
[Desktop Entry]
Exec=sh /home/mustafa/.0/startkit/start.sh
Icon=start
Name=StartMYBASH
Path=
Terminal=False
Type=Application
